# systemd/user/pass-pull.timer
[Unit]
Description=Pass pull git updates

[Timer]
OnBootSec=3min
OnUnitActiveSec=12h

[Install]
WantedBy=timers.target
